# ansible/roles/system/files/rc.local
#!/bin/sh -e

# Make sure the SSH host keys have been generated
ssh-keygen -A

exit 0
