# ansible/roles/system/files/rc.local
#!/bin/sh -e

# Make sure the SSH host keys have been generated
ssh-keygen -A

# Solves issue with black square in the upper left-hand corner
/bin/chvt 8 ; sleep 0.1 ; /bin/chvt 7

# Disable wifi power management
if [ "$(ip address | grep wlan)" ]; then
  iwconfig wlan0 power off
fi

exit 0
