# ansible/roles/network/templates/wifi-connect.service
[Unit]
Description=Anthias WiFi Connect
Documentation=https://github.com/Screenly/Anthias
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
User={{ anthias_user }}
WorkingDirectory=/home/{{ anthias_user }}/anthias
ExecStart=/home/{{ anthias_user }}/anthias/start_wifi_connect_service.sh

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anthias-wifi-connect

# Sandboxing
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX
LockPersonality=true
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
UMask=0027

[Install]
WantedBy=multi-user.target
