# ansible/roles/screenly/templates/anthias-host-agent.service
[Unit]
Description=Anthias Host Agent
Documentation=https://github.com/Screenly/Anthias
Requires=docker.service
After=docker.service
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Type=simple
User={{ anthias_user }}
WorkingDirectory=/home/{{ anthias_user }}/screenly
ExecStart=/home/{{ anthias_user }}/installer_venv/bin/python /home/{{ anthias_user }}/screenly/host_agent.py
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anthias-host-agent

# Sandboxing
# NoNewPrivileges, CapabilityBoundingSet and SystemCallFilter are
# *not* tightened here: host_agent invokes `sudo systemctl reboot`
# / `poweroff`, which requires setuid escalation, CAP_SYS_BOOT in
# the bounding set, and the reboot() syscall.
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
SystemCallArchitectures=native
UMask=0027

[Install]
WantedBy=multi-user.target
