# ansible/roles/screenly/templates/X.service
[Unit]
Description=X11
After=network.target

[Service]
ExecStart=/usr/bin/X
User=root
Restart=on-failure

[Install]
WantedBy=multi-user.target